# packages.d/hostapd.files
# Copyright (C) 2017 Hamish Coleman <hamish@zot.org>

rm      /etc/default/hostapd
rm      /etc/hostapd/ifupdown.sh
rm      /etc/init.d/hostapd
rm      /etc/network/if-post-down.d/hostapd
rm      /etc/network/if-pre-up.d/hostapd
keep    /usr/sbin/hostapd
keep    /usr/sbin/hostapd_cli

