# packages/cpufrequtils.files
# we dont have kernel modules inside the initramfs, so trying to load
# the modules with /etc/init.d/loadcpufreq just causes system errors

r   /etc/init.d/cpufrequtils
r   /etc/init.d/loadcpufreq
r   /etc/rc*.d/*cpufrequtils
r   /etc/rc*.d/*loadcpufreq
r   /usr/bin/cpufreq-aperf
k   /usr/bin/cpufreq-info
k   /usr/bin/cpufreq-set

# /usr/share/doc, /usr/share/locale and /usr/share/man
# are all dealt with centrally

