# packages/dpkg.files
rm  /etc/alternatives/README
rm  /etc/cron.daily/dpkg
k   /etc/dpkg/dpkg.cfg
rm  /etc/logrotate.d/dpkg
todo    /sbin/start-stop-daemon
k   /usr/bin/dpkg
k   /usr/bin/dpkg-deb
k   /usr/bin/dpkg-divert
k   /usr/bin/dpkg-maintscript-helper
k   /usr/bin/dpkg-query
k   /usr/bin/dpkg-split
k   /usr/bin/dpkg-statoverride
k   /usr/bin/dpkg-trigger
k   /usr/bin/update-alternatives
k   /usr/sbin/dpkg-divert
k   /usr/sbin/dpkg-statoverride
k   /usr/sbin/update-alternatives
todo    /usr/share/dpkg/abitable
todo    /usr/share/dpkg/cputable
todo    /usr/share/dpkg/ostable
todo    /usr/share/dpkg/triplettable
rm  /var/lib/dpkg/*-old
# rm   /var/lib/dpkg/info/*.md5sums # useful to keep so we can check conffiles
rm  /var/lib/dpkg/info/*.symbols
rm  /var/lib/dpkg/info/*.templates
rm  /var/log/dpkg.log
rm  /var/log/alternatives.log

# TODO
# - /etc/dpkg/dpkg.cfg can be smaller or eliminated all together

# /usr/share/doc, /usr/share/locale and /usr/share/man
# are all dealt with centrally

