# packages/ipmitool.files
rm  /etc/default/ipmievd
rm  /etc/init.d/ipmievd
rm  /lib/systemd/system/ipmievd.service
keep    /usr/bin/ipmitool
rm  /usr/lib/modules-load.d/ipmievd.conf
rm  /usr/sbin/ipmievd
rm  /usr/share/ipmitool/oem_ibm_sel_map

# /usr/share/doc, /usr/share/info, /usr/share/man and /usr/share/lintian
# are all dealt with centrally

