# packages/mactelnet-server.files
# Copyright (C) 2017 Hamish Coleman <hamish@zot.org>

keep    /etc/default/mactelnet-server
keep    /etc/init.d/mactelnet-server
keep    /etc/mactelnetd.users
keep    /usr/sbin/mactelnetd
# TODO
# - /etc/default/mactelnet-server can be replaced with a smaller one
# - /etc/init.d/mactelnet-server should eventually be replaced entirely with
#   a systemd unit file

# /usr/share/doc, /usr/share/locale and /usr/share/man
# are all dealt with centrally

