# packages/partimage-server.files
keep    /etc/default/partimaged
keep    /etc/init.d/partimaged
keep    /etc/pam.d/partimaged
keep    /etc/partimaged/partimage-certs.cnf
keep    /etc/partimaged/partimagedusers
keep    /usr/sbin/partimaged
keep    /usr/sbin/partimaged-passwd
rm,nopkg /var/log/partimage/partimaged.log

# /usr/share/doc, /usr/share/locale and /usr/share/man
# are all dealt with centrally

