# packages/util-linux.files
keep    /bin/dmesg
keep    /bin/lsblk
keep    /bin/more
rm  /bin/tailf
rm  /bin/wdctl
rm  /etc/default/hwclock
keep    /etc/init.d/hwclock.sh
keep    /etc/pam.d/runuser
keep    /etc/pam.d/runuser-l
keep    /lib/udev/hwclock-set
keep    /lib/udev/rules.d/85-hwclock.rules
keep    /sbin/agetty
keep    /sbin/blkdiscard
keep    /sbin/blkid
keep    /sbin/blockdev
keep    /sbin/cfdisk
rm  /sbin/chcpu
rm  /sbin/ctrlaltdel
keep    /sbin/fdisk
keep    /sbin/findfs
keep    /sbin/fsck
rm  /sbin/fsck.cramfs
rm  /sbin/fsck.minix
keep    /sbin/fsfreeze
keep    /sbin/fstrim
keep    /sbin/getty
keep    /sbin/hwclock
rm  /sbin/isosize
keep    /sbin/mkfs
rm  /sbin/mkfs.bfs
rm  /sbin/mkfs.cramfs
rm  /sbin/mkfs.minix
keep    /sbin/mkswap
keep    /sbin/pivot_root
rm  /sbin/raw
keep    /sbin/runuser
keep    /sbin/sfdisk
keep    /sbin/swaplabel
keep    /sbin/switch_root
rm  /sbin/wipefs
rm  /usr/bin/addpart
rm  /usr/bin/chrt
rm  /usr/bin/delpart
keep    /usr/bin/fallocate
keep    /usr/bin/flock
keep    /usr/bin/getopt
todo   /usr/bin/${ARCH} # TODO - will change for different arch
keep    /usr/bin/ionice
rm  /usr/bin/ipcmk
rm  /usr/bin/ipcrm
rm  /usr/bin/ipcs
rm  /usr/bin/line
keep    /usr/bin/linux32
keep    /usr/bin/linux64
keep    /usr/bin/lscpu
keep    /usr/bin/lslocks
rm  /usr/bin/mcookie
keep    /usr/bin/namei
keep    /usr/bin/nsenter
keep    /usr/bin/partx
rm  /usr/bin/pg
keep    /usr/bin/prlimit
rm  /usr/bin/rename.ul
rm  /usr/bin/resizepart
rm  /usr/bin/rev
keep    /usr/bin/setarch
rm  /usr/bin/setsid
keep    /usr/bin/setterm
rm  /usr/bin/taskset
keep    /usr/bin/unshare
keep    /usr/bin/utmpdump
rm  /usr/bin/whereis
#   /usr/bin/x86_64 - TODO - arch specific
rm  /usr/lib/mime/packages/util-linux
keep    /usr/sbin/fdformat
keep    /usr/sbin/ldattach
rm  /usr/sbin/readprofile
rm  /usr/sbin/rtcwake
rm  /usr/sbin/tunelp
rm  /usr/share/bash-completion/completions/addpart
keep    /usr/share/bash-completion/completions/blkdiscard
keep    /usr/share/bash-completion/completions/blkid
keep    /usr/share/bash-completion/completions/blockdev
keep    /usr/share/bash-completion/completions/cfdisk
rm  /usr/share/bash-completion/completions/chcpu
rm  /usr/share/bash-completion/completions/chrt
rm  /usr/share/bash-completion/completions/ctrlaltdel
rm  /usr/share/bash-completion/completions/delpart
keep    /usr/share/bash-completion/completions/fallocate
keep    /usr/share/bash-completion/completions/fdformat
keep    /usr/share/bash-completion/completions/fdisk
keep    /usr/share/bash-completion/completions/flock
keep    /usr/share/bash-completion/completions/fsck
rm  /usr/share/bash-completion/completions/fsck.cramfs
rm  /usr/share/bash-completion/completions/fsck.minix
keep    /usr/share/bash-completion/completions/fsfreeze
keep    /usr/share/bash-completion/completions/fstrim
keep    /usr/share/bash-completion/completions/getopt
keep    /usr/share/bash-completion/completions/hwclock
keep    /usr/share/bash-completion/completions/ionice
rm  /usr/share/bash-completion/completions/ipcrm
rm  /usr/share/bash-completion/completions/ipcs
rm  /usr/share/bash-completion/completions/isosize
keep    /usr/share/bash-completion/completions/ldattach
keep    /usr/share/bash-completion/completions/lsblk
keep    /usr/share/bash-completion/completions/lscpu
keep    /usr/share/bash-completion/completions/lslocks
rm  /usr/share/bash-completion/completions/mcookie
keep    /usr/share/bash-completion/completions/mkfs
keep    /usr/share/bash-completion/completions/mkfs.bfs
rm  /usr/share/bash-completion/completions/mkfs.cramfs
rm  /usr/share/bash-completion/completions/mkfs.minix
keep    /usr/share/bash-completion/completions/mkswap
keep    /usr/share/bash-completion/completions/more
keep    /usr/share/bash-completion/completions/namei
keep    /usr/share/bash-completion/completions/nsenter
keep    /usr/share/bash-completion/completions/partx
rm  /usr/share/bash-completion/completions/pg
keep    /usr/share/bash-completion/completions/pivot_root
keep    /usr/share/bash-completion/completions/prlimit
rm  /usr/share/bash-completion/completions/raw
rm  /usr/share/bash-completion/completions/readprofile
rm  /usr/share/bash-completion/completions/resizepart
rm  /usr/share/bash-completion/completions/rev
rm  /usr/share/bash-completion/completions/rtcwake
keep    /usr/share/bash-completion/completions/setarch
rm  /usr/share/bash-completion/completions/setsid
keep    /usr/share/bash-completion/completions/setterm
keep    /usr/share/bash-completion/completions/sfdisk
keep    /usr/share/bash-completion/completions/swaplabel
rm  /usr/share/bash-completion/completions/tailf
rm  /usr/share/bash-completion/completions/taskset
rm  /usr/share/bash-completion/completions/tunelp
keep    /usr/share/bash-completion/completions/unshare
keep    /usr/share/bash-completion/completions/utmpdump
rm  /usr/share/bash-completion/completions/wdctl
keep    /usr/share/bash-completion/completions/whereis
keep    /usr/share/bash-completion/completions/wipefs


# /usr/share/doc, /usr/share/info, /usr/share/man and /usr/share/lintian
# are all dealt with centrally

