# packages/util-linux.files
k   /bin/dmesg
k   /bin/lsblk
k   /bin/more
r   /bin/tailf
r   /bin/wdctl
r   /etc/default/hwclock
k   /etc/init.d/hwclock.sh
k   /etc/pam.d/runuser
k   /etc/pam.d/runuser-l
k   /lib/udev/hwclock-set
k   /lib/udev/rules.d/85-hwclock.rules
k   /sbin/agetty
k   /sbin/blkdiscard
k   /sbin/blkid
k   /sbin/blockdev
k   /sbin/cfdisk
r   /sbin/chcpu
r   /sbin/ctrlaltdel
k   /sbin/fdisk
k   /sbin/findfs
k   /sbin/fsck
r   /sbin/fsck.cramfs
r   /sbin/fsck.minix
k   /sbin/fsfreeze
k   /sbin/fstrim
k   /sbin/getty
k   /sbin/hwclock
r   /sbin/isosize
k   /sbin/mkfs
r   /sbin/mkfs.bfs
r   /sbin/mkfs.cramfs
r   /sbin/mkfs.minix
k   /sbin/mkswap
k   /sbin/pivot_root
r   /sbin/raw
k   /sbin/runuser
k   /sbin/sfdisk
k   /sbin/swaplabel
k   /sbin/switch_root
r   /sbin/wipefs
r   /usr/bin/addpart
r   /usr/bin/chrt
r   /usr/bin/delpart
k   /usr/bin/fallocate
k   /usr/bin/flock
k   /usr/bin/getopt
todo   /usr/bin/${ARCH} # TODO - will change for different arch
k   /usr/bin/ionice
r   /usr/bin/ipcmk
r   /usr/bin/ipcrm
r   /usr/bin/ipcs
r   /usr/bin/line
k   /usr/bin/linux32
k   /usr/bin/linux64
k   /usr/bin/lscpu
k   /usr/bin/lslocks
r   /usr/bin/mcookie
k   /usr/bin/namei
k   /usr/bin/nsenter
k   /usr/bin/partx
r   /usr/bin/pg
k   /usr/bin/prlimit
r   /usr/bin/rename.ul
r   /usr/bin/resizepart
r   /usr/bin/rev
k   /usr/bin/setarch
r   /usr/bin/setsid
k   /usr/bin/setterm
r   /usr/bin/taskset
k   /usr/bin/unshare
k   /usr/bin/utmpdump
r   /usr/bin/whereis
#   /usr/bin/x86_64 - TODO - arch specific
r   /usr/lib/mime/packages/util-linux
k   /usr/sbin/fdformat
k   /usr/sbin/ldattach
r   /usr/sbin/readprofile
r   /usr/sbin/rtcwake
r   /usr/sbin/tunelp
r   /usr/share/bash-completion/completions/addpart
k   /usr/share/bash-completion/completions/blkdiscard
k   /usr/share/bash-completion/completions/blkid
k   /usr/share/bash-completion/completions/blockdev
k   /usr/share/bash-completion/completions/cfdisk
r   /usr/share/bash-completion/completions/chcpu
r   /usr/share/bash-completion/completions/chrt
r   /usr/share/bash-completion/completions/ctrlaltdel
r   /usr/share/bash-completion/completions/delpart
k   /usr/share/bash-completion/completions/fallocate
k   /usr/share/bash-completion/completions/fdformat
k   /usr/share/bash-completion/completions/fdisk
k   /usr/share/bash-completion/completions/flock
k   /usr/share/bash-completion/completions/fsck
r   /usr/share/bash-completion/completions/fsck.cramfs
r   /usr/share/bash-completion/completions/fsck.minix
k   /usr/share/bash-completion/completions/fsfreeze
k   /usr/share/bash-completion/completions/fstrim
k   /usr/share/bash-completion/completions/getopt
k   /usr/share/bash-completion/completions/hwclock
k   /usr/share/bash-completion/completions/ionice
r   /usr/share/bash-completion/completions/ipcrm
r   /usr/share/bash-completion/completions/ipcs
r   /usr/share/bash-completion/completions/isosize
k   /usr/share/bash-completion/completions/ldattach
k   /usr/share/bash-completion/completions/lsblk
k   /usr/share/bash-completion/completions/lscpu
k   /usr/share/bash-completion/completions/lslocks
r   /usr/share/bash-completion/completions/mcookie
k   /usr/share/bash-completion/completions/mkfs
k   /usr/share/bash-completion/completions/mkfs.bfs
r   /usr/share/bash-completion/completions/mkfs.cramfs
r   /usr/share/bash-completion/completions/mkfs.minix
k   /usr/share/bash-completion/completions/mkswap
k   /usr/share/bash-completion/completions/more
k   /usr/share/bash-completion/completions/namei
k   /usr/share/bash-completion/completions/nsenter
k   /usr/share/bash-completion/completions/partx
r   /usr/share/bash-completion/completions/pg
k   /usr/share/bash-completion/completions/pivot_root
k   /usr/share/bash-completion/completions/prlimit
r   /usr/share/bash-completion/completions/raw
r   /usr/share/bash-completion/completions/readprofile
r   /usr/share/bash-completion/completions/resizepart
r   /usr/share/bash-completion/completions/rev
r   /usr/share/bash-completion/completions/rtcwake
k   /usr/share/bash-completion/completions/setarch
r   /usr/share/bash-completion/completions/setsid
k   /usr/share/bash-completion/completions/setterm
k   /usr/share/bash-completion/completions/sfdisk
k   /usr/share/bash-completion/completions/swaplabel
r   /usr/share/bash-completion/completions/tailf
r   /usr/share/bash-completion/completions/taskset
r   /usr/share/bash-completion/completions/tunelp
k   /usr/share/bash-completion/completions/unshare
k   /usr/share/bash-completion/completions/utmpdump
r   /usr/share/bash-completion/completions/wdctl
k   /usr/share/bash-completion/completions/whereis
k   /usr/share/bash-completion/completions/wipefs


# /usr/share/doc, /usr/share/info, /usr/share/man and /usr/share/lintian
# are all dealt with centrally

